# debian/picam-hailo/DEBIAN/prerm
#!/bin/sh
set -e
systemctl stop    picam-hailo 2>/dev/null || true
systemctl disable picam-hailo 2>/dev/null || true
